# config/systemd/yap-pipeline.timer
[Unit]
Description=Run YAP pipeline every hour (America/Bogota server time)

[Timer]
OnCalendar=hourly
Persistent=true
Unit=yap-pipeline.service

[Install]
WantedBy=timers.target
